# debian/etc/systemd/system/bytefense-watch.service
[Unit]
Description=Bytefense OS Security Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/bytefense/bin/bytefense-watch
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target