# BytefenseOS.AppDir/usr/share/bytefense/system/bytefense-dashboard-secure.service
[Unit]
Description=Bytefense OS Dashboard Seguro
After=network-online.target
Wants=network-online.target
Requires=bytefense-watch.service

[Service]
Type=simple
User=bytefense
Group=bytefense
WorkingDirectory=/opt/bytefense
ExecStart=/usr/bin/python3 /opt/bytefense/bin/bytefense-api-secure.py
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=10

# Seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bytefense /var/log
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Variables de entorno
Environment=PYTHONPATH=/opt/bytefense
Environment=BYTEFENSE_ENV=production
Environment=BYTEFENSE_LOG_LEVEL=INFO

# Límites de recursos
LimitNOFILE=1024
LimitNPROC=512
MemoryMax=256M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bytefense-dashboard

[Install]
WantedBy=multi-user.target