# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Post-installation script for Bytefense OS
echo "Configurando Bytefense OS..."

# Set proper permissions
chmod +x /opt/bytefense/bin/*
chown -R root:root /opt/bytefense

# Create bytefense user if it doesn't exist
if ! id "bytefense" &>/dev/null; then
    useradd -r -s /bin/false -d /opt/bytefense bytefense
fi

# Create necessary directories
mkdir -p /var/log/bytefense
mkdir -p /var/lib/bytefense
chown bytefense:bytefense /var/log/bytefense
chown bytefense:bytefense /var/lib/bytefense

# Enable and start services
systemctl daemon-reload
for service in bytefense-watch bytefense-dashboard bytefense-intel-updater; do
    if [ -f "/etc/systemd/system/.service" ]; then
        echo "Habilitando servicio: "
        systemctl enable 
        systemctl start 
    fi
done

echo "Bytefense OS instalado exitosamente"
echo "Dashboard disponible en: http://localhost:8080"

exit 0
