# debian/DEBIAN/prerm
#!/bin/bash
# Bytefense OS - Script de Pre-eliminación

set -e

log() {
    echo "[$(date '+%Y-%m-%d %H:%M:%S')] $1" | tee -a "/var/log/bytefense-remove.log"
}

case "$1" in
    remove|upgrade|deconfigure)
        log "🛑 Deteniendo servicios de Bytefense OS..."
        
        # Detener servicios
        systemctl stop bytefense-dashboard.service 2>/dev/null || true
        systemctl stop bytefense-watch.service 2>/dev/null || true
        
        # Deshabilitar servicios
        systemctl disable bytefense-dashboard.service 2>/dev/null || true
        systemctl disable bytefense-watch.service 2>/dev/null || true
        
        log "✅ Servicios detenidos"
        ;;
        
    failed-upgrade)
        log "❌ Fallo en actualización"
        ;;
        
    *)
        log "❓ Acción desconocida: $1"
        ;;
esac

exit 0